# templates/alertmanager.service.j2
[Unit]
Description=Alertmanager Service

[Service]
User=root
ExecStart={{ alertmanager_binary_path }}/alertmanager \
          --config.file=/etc/alertmanager/alertmanager.yml \
          --storage.path=/var/lib/alertmanager 

Restart=on-failure

[Install]
WantedBy=multi-user.target
